# systemd/assistant-mithu.service
[Unit]
Description=Tintu Assistant Bot — Mithu
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=assistant
WorkingDirectory=/opt/assistant-mithu
ExecStart=/opt/assistant-mithu/venv/bin/python -m app.main
Restart=on-failure
RestartSec=10
EnvironmentFile=/opt/assistant-mithu/.env

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
